# dropbox@.service
.include /usr/lib/systemd/system/dropbox@.service
[Service]
Environment=DISPLAY=:0
